# supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info

# Redis 消息队列
[program:redis]
command=redis-server --bind 0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

# FastAPI
[program:fastapi]
command=uvicorn backend:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

# Celery Worker
[program:celery_worker]
command=celery -A celery_backend worker --loglevel=info
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

# Flower 监控
[program:flower]
command=celery -A celery_backend flower --port=5555
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr